# configs/php_fpm/conf/laravel-queue.ini
[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/back/artisan queue:work --sleep=1
autostart=true
autorestart=true
user=local
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/back/queue.log
